# third_party/ros_x86_64/etc/catkin/profile.d/15.rosbash.zsh
# generated from rosbash/env-hooks/15.rosbash.zsh.em

if [ -z "$CATKIN_ENV_HOOK_WORKSPACE" ]; then
  CATKIN_ENV_HOOK_WORKSPACE="/apollo/third_party/ros_x86_64"
fi
. "$CATKIN_ENV_HOOK_WORKSPACE/share/rosbash/roszsh"
